# roles/solana_validator_bootstrap/templates/solana-sys-tuner.service.j2
[Unit]
Description=Solana System Tuner Service
After=network.target syslog.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart={{ solana_home }}/.local/share/solana/install/active_release/bin/solana-sys-tuner --user {{ solana_user }}

[Install]
WantedBy=multi-user.target
